# deploy/systemd/aimax-report-auto-guidance.service
[Unit]
Description=AIMAX user-facing report guidance updater

[Service]
Type=oneshot
WorkingDirectory=/home/ubuntu/.openclaw/workspace/AIMAX-AI-Staff-Management
ExecStart=/usr/bin/python3 /home/ubuntu/.openclaw/workspace/AIMAX-AI-Staff-Management/scripts/aimax_report_auto_guidance.py --min-age-minutes 15 --lookback-days 14
