# deploy/systemd/aimax-error-repair-agent.timer
[Unit]
Description=Run AIMAX error repair Codex agent every 30 minutes

[Timer]
OnBootSec=6min
OnActiveSec=3min
OnUnitActiveSec=30min
Persistent=true
Unit=aimax-error-repair-agent.service

[Install]
WantedBy=timers.target
