# deploy/systemd/aimax-report-auto-guidance.timer
[Unit]
Description=Run AIMAX user-facing report guidance updater every 5 minutes

[Timer]
OnBootSec=4min
OnActiveSec=2min
OnUnitActiveSec=5min
Persistent=true
Unit=aimax-report-auto-guidance.service

[Install]
WantedBy=timers.target
